# supervisord.conf
[supervisord]
loglevel=debug
nodaemon=true

[program:scheduler]
command=java -Xmx1g -jar target/scheduler.jar -file ../config/stucco.yml -section scheduler
autostart=true
autorestart=true
